# Library/org.zfs.snapshot.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>org.zfs.snapshot.weekly</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/sbin/zfs-auto-snapshot.sh</string>
			<string>--syslog</string>
			<string>--keep=8</string>
			<string>--label=weekly</string>
			<string>--send-incr=zones/tank/archive/mac/nicks-macbook-air</string>
			<string>--remove-local=4</string>
			<string>--create</string>
			<string>--fallback</string>
			<string>--destroy</string>
			<string>//</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
		<key>StartCalendarInterval</key>
		<dict>
			<key>Minute</key>
			<integer>22</integer>
			<key>Hour</key>
			<integer>23</integer>
			<key>Weekday</key>
			<integer>0</integer>
			
		</dict>
	</dict>
</plist>
